# files/service/mysql/run
#!/bin/bash
mysqld
# --datadir=/opt/mysql/mysql/data