# files/service/mysql/run
#!/bin/bash
/sbin/setuser mysql mysqld
# --datadir=/opt/mysql/mysql/data