# files/service/apache2/run
#!/bin/bash
apachectl --with-mysqli -DFOREGROUND